# ssl/renew.crontab
# m h dom mon dow command
30 2 * * 1 sudo /usr/bin/letsencrypt renew
